# system/etc/systemd/system/mpd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/mpd --no-daemon /home/hug/.config/mpd/mpd.conf
User=hug
Group=hug